# conf/play-watcher.path
[Path]
PathModified=/home/marc/play/main

[Install]
WantedBy=multi-user.target
