# conf/play-go-watcher.path
[Path]
PathModified=/home/marc/play-go/main

[Install]
WantedBy=multi-user.target
