# conf/play-go.path
[Unit]
Description=Watch /home/marc/play-go for changes

[Path]
PathModified=/home/marc/play-go/main

[Install]
WantedBy=multi-user.target
